# supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:concierge]
command=npm start
directory=/var/app
autostart=true
autorestart=true
stopsignal = QUIT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
